# bin/Debug/net6.0/identity.autofarmer.net.service
[Unit]
Description=checkinfo.consumer
ConditionPathExists=/root/Identity/output/identity/
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/Identity/output/identity/
ExecStart=/root/Identity/output/identity/Identity
LimitNOFILE=102400
LimitNOFILESoft=102400
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target